# deploy/lattice-relay.service
[Unit]
Description=Lattice Relay Hub
Documentation=https://github.com/djsph001/lattice-node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lattice
Group=lattice
ExecStart=/opt/lattice-node/target/release/lattice-node \
    --name relay-hub \
    --port 4001 \
    --listen-addr 0.0.0.0 \
    --identity-dir /opt/lattice-node/data/identity \
    --storage-dir /opt/lattice-node/data/storage \
    --relay-server
Restart=always
RestartSec=5
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/lattice-node/data
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lattice-relay

[Install]
WantedBy=multi-user.target
